# src/etc/services.d/01-php-fpm/finish
#!/bin/sh -e

echo "[services.d] stopping service: php5-fpm"
