# src/etc/services.d/01-php-fpm/run
#!/bin/sh -e

touch /mnt/log/php-error.log
chown nobody:nobody /mnt/log/php-error.log

echo "[php5-fpm] start"
php-fpm -F
