# src/etc/services.d/01-php-fpm/run
#!/bin/sh -e

touch /mnt/log/php-error.log
chown nobody:nobody /mnt/log/php-error.log

echo "[services.d] starting service: php5-fpm"
php-fpm -F
